# bin/kiosk.service
# save file to /lib/systemd/system/
# then sudo systemctl enable kiosk.service
[Unit]
Description=Chromium Kiosk
Wants=graphical.target
After=graphical.target

[Service]
Environment=DISPLAY=:0.0
Environment=XAUTHORITY=/home/john/.Xauthority
Type=simple
ExecStart=/bin/bash /home/john/scoreboard/bin/start-kiosk.sh
Restart=on-abort
User=john
Group=john

[Install]
WantedBy=graphical.target